# supervisord.conf
[program:hhvm]
command=/usr/bin/hhvm --mode server -c /etc/hhvm/server.ini
autorestart=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autorestart=true

